# files/etc/systemd/system/mariadb-resume.service
[Unit]
Description=Local system resume actions
After=suspend.target

[Service]
Type=simple
ExecStart=/usr/bin/systemctl start mariadb.service

[Install]
WantedBy=suspend.target
